# frontend/deploy-temp/portfoliopulse-frontend.service
[Unit]
Description=PortfolioPulse Frontend Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/portfoliopulse-frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfoliopulse-frontend

[Install]
WantedBy=multi-user.target
